# dist/fix-ip-in-index.service
# /etc/systemd/system/fix-ip-in-index.service
# This service will automatically fix the ip address in the index.html file in the
# swagger-ui/dist area so that it matches the address of this host.
[Unit]
After=network-online.target

[Service]
Type=oneshot
ExecStart=/root/swagger-ui/dist/fix-ip-in-index.sh
StandardOutput=syslog
StandardError=syslog


[Install]
WantedBy=multi-user.target
